# root/etc/services.d/dockerd/run
#!/usr/bin/with-contenv bash
exec /root/.dropbox-dist/dropboxd
